# screen/.config/profile.d/30-screen.sh
# Set screen configuration location
export SCREENRC="${CONFIG}/screenrc"
