# host-cactus-vm/.config/profile.d/99-host-cactus-vm.sh
export SSL_CERT_DIR="/etc/ssl/certs:${HOME}/.config/certs"
export SSH_ASKPASS="/usr/bin/x11-ssh-askpass"

export PICO_SDK_PATH="${HOME}/src/pico-sdk"
