# vim/.config/profile.d/30-vim.sh
export VIMINIT="source ${CONFIG}/vim/vimrc"
export EDITOR=vim
