# ruby/.config/profile.d/30-ruby.sh
# Install Ruby Gems to ~/gems
export GEM_HOME="${HOME}/gems"
export PATH="${PATH}:${HOME}/gems/bin"
