# rustlang/.config/profile.d/50-rustlang.sh
# Add path to rust dev tools
export PATH="${HOME}/.cargo/bin:${PATH}"
