# user-usr/.config/profile.d/50-user-usr.sh
# Add paths in ~/usr
export PATH="${HOME}/usr/bin:${HOME}/usr/sbin:${PATH}"
export LD_LIBRARY_PATH="${HOME}/usr/lib"
export MANPATH="${HOME}/usr/man:${HOME}/usr/share/man:${MANPATH}"
export PKG_CONFIG_PATH="${HOME}/usr/lib/pkgconfig:${PKG_CONFIG_PATH}"
