# host-ATOMICUS-M-Q1FM/.config/profile.d/99-gpg-ssh-agent.sh
# Use gpg-agent for SSH authentication
export SSH_AUTH_SOCK="${HOME}/.gnupg/S.gpg-agent.ssh"
